# files/system/common/etc/profile.d/ananya-fastfetch.sh
alias fastfetch="/usr/libexec/ananya-fastfetch"
